# airpup
### BEGIN INIT INFO
# Provides: airpup
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start AirPup at boot time
# Description: Start AirPup at boot time
### END INIT INFO

#! /bin/sh

PATH=/bin:/usr/bin:/sbin:/usr/sbin
DAEMON=/home/pi/airpup/airpup.sh
PIDFILE=/var/run/airpup.pid

test -x $DAEMON || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
     log_daemon_msg "Starting Airpup"
     start_daemon -p $PIDFILE $DAEMON
     log_end_msg $?
   ;;
  stop)
     log_daemon_msg "Stopping Airpup"
     killproc -p $PIDFILE $DAEMON
     PID=`ps x | grep airpup | head -1 | awk '{print $1}'`
     kill -9 $PID
     log_end_msg $?
   ;;
  force-reload|restart)
     $0 stop
     $0 start
   ;;
  status)
     status_of_proc -p $PIDFILE $DAEMON airpup && exit 0 || exit $?
   ;;
 *)
   echo "Usage: /etc/init.d/airpup {start|stop|restart|force-reload|status}"
   exit 1
  ;;
esac

exit 0